# scripts/openclaw-autostart.service
[Unit]
Description=OpenClaw Auto-Start Service
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash /root/.openclaw/workspace/scripts/startup.sh
User=root
WorkingDirectory=/root/.openclaw/workspace
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/root/.openclaw/.env

[Install]
WantedBy=multi-user.target
